# scripts/scheduling/rag-sync.timer
[Unit]
Description=Run Dreamweaving RAG sync daily at 3 AM

[Timer]
# Run at 3:00 AM every day
OnCalendar=*-*-* 03:00:00

# If the system was off when timer should have fired, run immediately
Persistent=true

# Randomize start time by up to 5 minutes to avoid thundering herd
RandomizedDelaySec=300

# Accuracy: allow firing up to 1 minute late for efficiency
AccuracySec=1min

[Install]
WantedBy=timers.target
